# gitea/rootfs/etc/services.d/gitea/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the Gitea service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here
set -e

declare CONFIG_PATH=/data/options.json
bashio::log.green "[Info] Starting addon process gitea..."

for FOLDER in /data/gitea/conf /data/gitea/log /data/git /data/ssh; do
    if [ ! -d ${FOLDER} ]; then
        mkdir -p ${FOLDER}
    fi
done

if [ ! -d /data/git/.ssh ]; then
 mkdir -p /data/git/.ssh
 chmod 700 /data/git/.ssh
fi

## Run your program
exec su-exec $USER /app/gitea/gitea web

