# gitea/rootfs/etc/services.d/gitea/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the Gitea service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here
set -e

declare CONFIG_PATH=/data/options.json
bashio::log.green "[Info] Starting addon process gitea..."

for FOLDER in /data/gitea/conf /data/gitea/log /data/git /data/ssh; do
    if [ ! -d ${FOLDER} ]; then
        mkdir -p ${FOLDER}
    fi
done

if [ ! -d /data/git/.ssh ]; then
 mkdir -p /data/git/.ssh
 chmod 700 /data/git/.ssh
fi

# only chown if current owner is not already the gitea ${USER}. No recursive check to save time
if ! [[ $(ls -ld /data/gitea | awk '{print $3}') = ${USER} ]]; then chown -R ${USER}:git /data/gitea || true; fi
if ! [[ $(ls -ld /app/gitea | awk '{print $3}') = ${USER} ]]; then chown -R ${USER}:git /app/gitea || true; fi
if ! [[ $(ls -ld /data/git | awk '{print $3}') = ${USER} ]]; then chown -R ${USER}:git /data/git || true; fi

# Apply configuration through ENV variables
bashio::config.require domain
bashio::config.require.ssl

DOMAIN=$(bashio::config 'domain')
PORT=$(bashio::addon.port 3000)
if bashio::config.true 'ssl'; then
  bashio::log.info "SSL is enabled"
  PROTOCOL=https
  export GITEA__server__CERT_FILE="/ssl/fullchain.pem"
  export GITEA__server__KEY_FILE="/ssl/privkey.pem"
else
  PROTOCOL=http
if
if bashio::config.has_value 'root_url'; then
  ROOT_URL=$(bashio::config 'root_url')
  bashio::log.blue "ROOT_URL set, using value : ${ROOT_URL}"
else
  ROOT_URL="${PROTOCOL}://${DOMAIN}:${PORT}"
  bashio::log.blue "ROOT_URL not set, using extrapolated value : ${ROOT_URL}"
fi
export GITEA__server__PROTOCOL=${PROTOCOL}
export GITEA__server__ROOT_URL=${ROOT_URL}
export GITEA__server__HTTP_PORT=${PORT}

# Patch nginx configuration
bashio::config.require.use_nginx
if bashio.config.true "use_nginx"; then
  bashio::log.info "Configure Nginx reverse proxy"
  if ping -c 1 core-nginx-proxy >/dev/null 2>&1; then
    bashio::log.info "Current Nginx configuration"
    bashio::addon.options core_nginx_proxy
    echo

    bashio::log.info "Adding custom config to /share/nginx_proxy/nginx_gitea.conf"
    mkdir -p /share/nging_proxy
    sed "s/__DOMAIN__/${DOMAIN}/g; s/__GITEA_PORT__/${PORT}/g" /app/nginx_gitea.conf > /share/nginx_proxy/nginx_gitea.conf
    cat /share/nginx_proxy/nginx_gitea.conf
  else
    bashio::log.warning "Nginx is not running"
  fi
else
  bashio::log.waring "Nginx reverse proxy is disabled"
fi

## Run your program
exec su-exec $USER /app/gitea/gitea web

